# app/app/services/templates/services/service.conf
{% autoescape off %}

[program:{{ program }}]
autorestart={{ autorestart|default:True|yesno:'true,false' }}
autostart=false
command={{ command }}
stderr_events_enabled=true
stderr_logfile=/proc/1/fd/2
stderr_logfile_maxbytes=0
stdout_events_enabled=true
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
{% if extras %}
{% for key, value in extras.items %}
{{ key }}={{ value }}
{% endfor %}
{% endif %}

{% endautoescape %}
