# harbor/image/etc/supervisor/conf.d/icecast.conf
[program:nginx]
command = nginx -g 'daemon off;'
autorestart = true
redirect_stderr = true
stdout_logfile = /var/log/icecast
