# harbor/image/etc/supervisor/conf.d/icecast.conf
[program:icecast]
command = icecast2 -c /etc/icecast2/icecast.xml
user = icecast2
group = icecast
autorestart = true
redirect_stderr = true
stdout_logfile = /var/log/icecast.log
