# harbor/image/etc/supervisor/conf.d/upstream.conf
[program:upstream]
command = /mnt/run.sh upstream.liq
user = user
group = user
autorestart = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = /var/log/upstream.log
