# admin/app/carb/templates/services/service.conf
{% load carb %}
[program:{{ program }}]
autorestart={{ autorestart|boolstr }}
autostart={{ autostart|boolstr }}
command={{ command }}
stderr_events_enabled=true
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_events_enabled=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
{% if extras %}{% for key, value in extras.items %}{{ key }}={{ value }}
{% endfor %}{% endif %}
